# install/deb/vhestia-panel-watchdog.service
[Unit]
Description=VHestia Panel Watchdog Check

[Service]
Type=oneshot
ExecStart=/usr/local/hestia/bin/v-check-vhestia-panel
