# install/deb/vhestia-panel.service
[Unit]
Description=VHestia Panel (PM2 managed Node.js application)
Documentation=https://github.com/vhestiacp
After=network.target

[Service]
Type=forking
User=root
LimitNOFILE=infinity
LimitNPROC=infinity
LimitCORE=infinity
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/hestia/bin
Environment=PM2_HOME=/root/.pm2
PIDFile=/root/.pm2/pm2.pid
Restart=on-failure
RestartSec=5

ExecStart=/usr/bin/pm2 start /usr/local/hestia/web_v2/ecosystem.config.cjs --env production
ExecReload=/usr/bin/pm2 reload vhestia-panel --update-env
ExecStop=/usr/bin/pm2 stop vhestia-panel

[Install]
WantedBy=multi-user.target
