# install/deb/vhestia-panel.service
[Unit]
Description=VHestiaCP Web Panel v2
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/hestia/web_v2/server
ExecStart=/usr/local/bin/node src/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=8083
Environment=HESTIA=/usr/local/hestia
Environment=VHESTIA=/usr/local/hestia
Environment=SSL_KEY=/usr/local/hestia/ssl/certificate.key
Environment=SSL_CERT=/usr/local/hestia/ssl/certificate.crt

# Logging
StandardOutput=append:/usr/local/hestia/log/vhestia-panel.log
StandardError=append:/usr/local/hestia/log/vhestia-panel-error.log

# Security
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
